# containers/Deenish-Island/site/crontab
*/10 * * * * /usr/local/bin/python /root/site.py 
