# containers/El-Campello/waves/crontab
30 * * * * /usr/local/bin/python /root/rs.py 
