# multidocker/rs/crontab
08 10 * * * /usr/local/bin/python /root/rs.py 
