# containers/El-Campello/telebot/crontab
30 * * * * /usr/local/bin/python /root/main.py 
