# containers/Deenish-Island/rbn/crontab
0 22 * * * /usr/local/bin/python /root/main.py 
