# containers/Deenish-Island/iwbn/crontab
28 * * * * /usr/local/bin/python /root/iwbn.py 

