# containers/Deenish-Island/iwbn/crontab
35 * * * * /usr/local/bin/python /root/iwbn.py 

