# containers/Deenish-Island/chl/crontab
35 2-23/2 * * * /usr/local/bin/python /root/rs.py 
