# multidocker/rs-2/crontab
38 20 * * * /usr/local/bin/python /root/rs.py 
