# multidocker/model/crontab
35 16 * * * /usr/local/bin/python /root/model.py 
