# multidocker/model/crontab
0 8 * * * /usr/local/bin/python /root/model.py 
