# multidocker/site-2/crontab
1-59/10 * * * * /usr/local/bin/python /root/site.py 
