# multidocker/lptm/crontab
25 16 * * * /opt/conda/bin/python /root/lptm.py 
